# mydhcpd/root/etc/hotplug.d/iface/25-dhcpd
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

/etc/init.d/dhcpd enabled && /etc/init.d/dhcpd restart
